# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

# ─── PID 0/1 stack: virtual display ───────────────────────────────────────────
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24
priority=10
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
environment=DISPLAY=":99"

[program:fluxbox]
command=/usr/bin/fluxbox
priority=20
autostart=true
autorestart=true
startsecs=2
environment=DISPLAY=":99"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:x11vnc]
command=/usr/bin/x11vnc -display :99 -forever -nopw -rfbport 5900 -shared -quiet
priority=30
autostart=true
autorestart=true
startsecs=2
environment=DISPLAY=":99"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:novnc]
command=/usr/bin/websockify --web /usr/share/novnc 6080 localhost:5900
priority=40
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

# ─── Chrome via wrapper (clears SingletonLock, waits Xvfb) ────────────────────
[program:chrome]
command=/usr/local/bin/start-chrome.sh
priority=50
autostart=true
autorestart=true
startsecs=10
startretries=10
environment=DISPLAY=":99"
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

# ─── CDP forward 9223 → 9222 ──────────────────────────────────────────────────
[program:cdp-forward]
command=/usr/bin/socat TCP-LISTEN:9223,fork,reuseaddr TCP:localhost:9222
priority=60
autostart=true
autorestart=true
startsecs=2
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

# ─── search-proxy FastAPI on :3100 (localhost only) ───────────────────────────
[program:search-proxy]
command=/usr/bin/python3 /app/server.py
priority=70
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

# ─── SearxNG (granian) on :8080 ───────────────────────────────────────────────
[program:searxng]
command=/usr/local/bin/start-searxng.sh
priority=80
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

# ─── Business-level watchdog: detects "process alive but search broken" ───────
[program:watchdog]
command=/usr/local/bin/watchdog.sh
priority=90
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
